# deploy/templates/ldapi-server-configwatcher.path
[Unit]
Description=Restart LD API on config change

[Path]
PathChanged=/etc/ldapi-server/ldapi-server.ini

[Install]
WantedBy=multi-user.target
